# osm-tiles/supervisord.conf
[supervisord]
; supervisord log file
logfile=/var/log/supervisor/supervisord.log
; info, debug, warn, trace
loglevel=debug
; number of startup file descriptors
minfds=1024
; number of process descriptors
minprocs=200
; default user
user=root
; where child log files will live
childlogdir=/var/log/supervisor/
; run supervisord as a daemon
nodaemon=false

[unix_http_server]
file=/var/run/supervisor.sock ; (the path to the socket file)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL for a unix socket

[program:apache2]
command=/usr/sbin/apache2ctl -DFOREGROUND
environment=APACHE_LOG_DIR=/var/log/apache2
killasgroup=true
stopasgroup=true

[program:renderd]
user=www-data
command=sh -c "/usr/local/bin/datasource-settings.sh && /usr/local/bin/renderd -f -c /usr/local/etc/renderd.conf"
